# kuzco.service
[Unit]
Description=Kuzco AI Assistant CLI - Service to ensure availability
After=network.target

[Service]
Type=oneshot
User=dev
WorkingDirectory=/home/dev/Documents/GitHub/python/ai-cli
ExecStart=/bin/true
RemainAfterExit=yes
Environment=PATH=/home/dev/Documents/GitHub/python/ai-cli/venv/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
